# 99-10moons-tablet.rules
# 10moons tablet - permissões USB
SUBSYSTEM=="usb", ATTR{idVendor}=="08f2", ATTR{idProduct}=="6811", GROUP="plugdev", MODE="0664"
SUBSYSTEM=="input", ATTRS{idVendor}=="08f2", ATTRS{idProduct}=="6811", GROUP="plugdev", MODE="0664"

# Iniciar driver via systemd (vincula ao ciclo de vida do dispositivo)
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="08f2", ATTR{idProduct}=="6811", TAG+="systemd", ENV{SYSTEMD_WANTS}="10moons-tablet.service"
ACTION=="remove", SUBSYSTEM=="usb", ATTR{idVendor}=="08f2", ATTR{idProduct}=="6811", RUN+="/usr/bin/systemctl stop 10moons-tablet.service"
